# ac-server.service
[Unit]
Description=AC IR Server

[Service]
Type=simple
User=root
ExecStart=/usr/bin/ac-ir-cmd

[Install]
WantedBy=multi-user.target
