# ac-server.service
[Unit]
Description=AC IR Server

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/ac-ir-cmd

[Install]
WantedBy=multi-user.target
